# 2/debian11/2.4.3/etc/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisord]
logfile=redirect_stdout
loglevel=info
user=root

[program:php]
command=php-fpm -R
autorestart=true
autostart=true
priority=1

[program:cron]
command=cron -f
autostart=${ENABLE_CRONJOBS}
autorestart=true
priority=2

[program:nginx]
command=nginx -g "daemon off;"
autorestart=true
autostart=true
startsecs=30
priority=999
